# systemd/cryptotrader.service
[Unit]
Description=CryptoTrader Bitcoin Trading Bot
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/Bitcoiner
Environment="PATH=/home/YOUR_USERNAME/Bitcoiner/venv/bin:/usr/bin"
ExecStart=/home/YOUR_USERNAME/Bitcoiner/venv/bin/python main_trader.py --interval 5
Restart=always
RestartSec=10
StandardOutput=append:/home/YOUR_USERNAME/Bitcoiner/logs/cryptotrader.log
StandardError=append:/home/YOUR_USERNAME/Bitcoiner/logs/cryptotrader_error.log

# リソース制限
LimitNOFILE=65536
TimeoutStopSec=30

# クラッシュ時の自動再起動設定
StartLimitInterval=200
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
